# leonardo.service
[Unit]
Description=Leonardo — Sports Prediction Scheduler
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/hlnx4/leonardo
ExecStart=/home/hlnx4/anaconda3/bin/python3 /home/hlnx4/leonardo/scheduler.py
Restart=always
RestartSec=10
StandardOutput=append:/home/hlnx4/leonardo/leonardo.log
StandardError=append:/home/hlnx4/leonardo/leonardo.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
